# system/fifotv.service
[Unit]
Description=FIFOtv Smart TV v2
After=graphical.target

[Service]
Type=simple
User=tv
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/home/tv/smarttv
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
